# .config/systemd/user/udiskie.service
[Unit]
Description=udiskie automounting daemon

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
