# .crontab
#!/usr/bin/crontab
#
# ONCE A DAY ###################################################################
#
# sound a gong at 22:00 - Daily Show
0 22 * * 1-4 exec /home/jpkotta/bin/gong
55 14 * * 1,3,5 exec /home/jpkotta/bin/gong
55 12 * * 1,3,5 exec /home/jpkotta/bin/gong
0 0 * * * crontab -l > /home/jpkotta/.crontab
#
#
# ONCE A WEEK ##################################################################
#
# make all files in my home inaccessible to others
0 0 * * 0 exec chmod o-rwx /home/jpkotta/ > /dev/null
#
#
#
# TWICE A MONTH ################################################################
#
# back up all config files and some data files in ~
#0 0 1,15 * * exec /home/jpkotta/bin/backup -config > /dev/null
#0 0 1,15 * * exec /home/jpkotta/bin/backup -data > /dev/null
#
# clean up ~/backup/, delete all files over a month old
#10 0 1,15 * * exec /home/jpkotta/bin/trash `/bin/find /home/jpkotta/backup -mtime +30` > /dev/null
#
#
#
