# .etc/paccache.timer
[Unit]
Description=Clean the package cache

[Timer]
Persistent=false
OnCalendar=weekly
Unit=paccache.service

[Install]
WantedBy=default.target