# .etc/pacman-download.service
[Unit]
Description=Download packages
RefuseManualStop=true

[Service]
Type=oneshot
ExecStart=/usr/bin/checkupdates -d
