# .etc/pacman-download.timer
[Unit]
Description=Download packages

[Timer]
Persistent=false
OnCalendar=daily
Unit=pacman-download.service

[Install]
WantedBy=default.target